# tag-linux/config/systemd/user/gdrive-sync.service
[Unit]
Description=Sync git org file repo

[Service]
ExecStart=/usr/bin/zsh %h/bin/gdrive-sync